# ansible/roles/btrfs-backup/templates/prune.service.j2
[Unit]
Description=BTRFS Snapshot Prune

[Service]
Type=oneshot
ExecStart=/usr/local/bin/btrfs-prune.sh

