# ansible/roles/btrfs-backup/templates/send.service.j2
[Unit]
Description=BTRFS Send Backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/btrfs-send.sh

