# ansible/roles/btrfs-backup/templates/snapshot.service.j2
[Unit]
Description=BTRFS Snapshot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/btrfs-snapshot.sh

