# ansible/roles/btrfs-target-prune/templates/prune-remote.service.j2
[Unit]
Description=BTRFS Remote Prune

[Service]
Type=oneshot
ExecStart=/usr/local/bin/btrfs-prune-remote.sh

